# data/systemd/system/incus.service.d/kapsule-log-dir.conf
# SPDX-FileCopyrightText: 2024-2026 Frostyard
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Drop-in for incus.service:
# - Ensure log directory exists before starting
# - Use wheel group instead of incus-admin for easier access

[Service]
ExecStartPre=-/usr/bin/mkdir -p /var/log/incus
ExecStart=
ExecStart=/usr/lib/incus/incusd --group=wheel --logfile=/var/log/incus/incusd.log
