# data/systemd/system/incus-user.socket.d/kapsule-socket-mode.conf
# SPDX-FileCopyrightText: 2024-2026 KDE Community
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Drop-in for incus-user.socket to allow all users access

[Socket]
SocketMode=0666